# bot-trading/deploy/bot-trading.service
[Unit]
Description=Binance Spot Scalping Bot
After=network.target

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/opt/bot-trading
EnvironmentFile=/opt/bot-trading/.env
ExecStart=/usr/bin/java -jar /opt/bot-trading/bot-trading.jar
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
